# uptimeKuma/rootfs/etc/services.d/uptimekuma/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Uptime Kuma
# Starts Uptime Kuma
# ==============================================================================
export NODE_PATH=/opt/node_modules:/data/node_modules
declare -a options

cd /opt || bashio::exit.nok "Could not change directory to Uptime Kuma"

bashio::log.info "Starting Uptime Kuma..."

exec node server/server.js"
